# deploy/conduit.service
[Unit]
Description=Conduit Gateway
Documentation=https://github.com/conduit/conduit
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=conduit
Group=conduit
WorkingDirectory=/home/conduit/conduit

# Binary and config
ExecStart=/home/conduit/conduit/bin/conduit server --config /home/conduit/conduit/config.json

# Environment
Environment=PATH=/usr/local/go/bin:/usr/bin:/bin:/usr/local/bin
EnvironmentFile=-/home/conduit/conduit/.conduit-secrets.env
# Optional: override default data directory (~/.conduit/)
# Environment=CONDUIT_DATA_DIR=/home/conduit/.conduit

# Graceful shutdown (gateway handles SIGTERM)
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/conduit/
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=conduit

[Install]
WantedBy=multi-user.target
